# src/dreamrtspserver.service
[Unit]
Description=Dreambox Real-Time Streaming Protocol Server
Requires=dbus.socket dev-venc0.device
After=dbus.socket dev-venc0.device

[Service]
Type=dbus
BusName=com.dreambox.RTSPserver
Restart=on-failure
ExecStart=/usr/bin/dreamrtspserver
StandardOutput=null

[Install]
WantedBy=multi-user.target
